# Automation/crontab




*/1 * * * * /home/pi/connectssh.sh >/dev/null 2>&1
@reboot /home/pi/startsakis.sh >/dev/null 2>&1
@reboot /home/pi/startsensors.sh >/dev/null 2>&1
